# deploy/ansible/roles-os/1.17-generic-pacemaker/templates/pacemaker.timer.j2
[Unit]
Description=Delay start of pacemaker.service after boot

[Timer]
OnBootSec={{ pacemaker_start_delay }}
Unit=pacemaker.service

[Install]
WantedBy=timers.target
